# beanstalkd/usr/lib/systemd/system/beanstalkd.service.d/beanstalkd.conf
[Service]
Environment=BEANSTALKD_ADDR=0.0.0.0
Environment=BEANSTALKD_PORT=11300
Environment=BEANSTALKD_USER=beanstalkd
Environment=BEANSTALKD_MAX_JOB_SIZE=65535
Environment=BEANSTALKD_BINLOG_DIR=/var/lib/beanstalkd/binlog
Environment=BEANSTALKD_BINLOG_SIZE=104857600
Environment=BEANSTALKD_BINLOG_FSYNC_PERIOD=0
